# gajim.service
[Unit]
Description=Gajim XMPP Client

[Service]
ExecStart=/usr/bin/gajim
Restart=on-abort

[Install]
WantedBy=graphical.target
