# kwin.service
[Unit]
Description=KWin Window Manager

[Service]
ExecStart=/usr/bin/kwin
Restart=on-abort

[Install]
WantedBy=graphical.target
Alias=wm.service
