# pidgin.service
[Unit]
Description=Pidgin Messenger

[Service]
ExecStart=/usr/bin/pidgin
Restart=on-abort

[Install]
WantedBy=graphical.target
