# plasma.service
[Unit]
Description=Plasma Desktop

[Service]
ExecStart=/usr/bin/plasma-desktop --nofork
Restart=on-abort

[Install]
WantedBy=graphical.target
