# qbittorrent.service
[Unit]
Description=qBittorrent

[Service]
ExecStart=/usr/bin/qbittorrent
Restart=on-abort

[Install]
WantedBy=graphical.target
