# thunderbird.service
[Unit]
Description=Thunderbird Mail Client

[Service]
ExecStart=/usr/bin/thunderbird
Restart=on-abort

[Install]
WantedBy=graphical.target
