# vlc.service
[Unit]
Description=VLC media player.

[Service]
ExecStart=/usr/bin/vlc -f
Restart=on-abort

[Install]
WantedBy=graphical.target
