# xfce4-panel.service
[Unit]
Description=Xfce4 Panel

[Service]
ExecStart=/usr/bin/xfce4-panel
Restart=on-abort

[Install]
WantedBy=graphical.target
