# xfdesktop.service
[Unit]
Description=Xfdesktop

[Service]
ExecStart=/usr/bin/xfdesktop
Restart=on-abort

[Install]
WantedBy=graphical.target
