# xfwm4.service
[Unit]
Description=Xfce4 Window Manager

[Service]
ExecStart=/usr/bin/xfwm4
Restart=on-abort

[Install]
WantedBy=graphical.target
