# xscreensaver.service
[Unit]
Description=XScreenSaver

[Service]
ExecStart=/usr/bin/xscreensaver -nosplash
Restart=on-abort

[Install]
WantedBy=graphical.target
